# roles/psdash/templates/psdash-agent-supervisor.conf.j2
[program:psdash-agent]
command={{ psdash.path }}/bin/psdash -a -p 5000 -l /var/log/*.log -l /var/log/**/*.log
;environment=PATH="{{ psdash.path }}/bin:{{ lookup('env', 'PATH') }}"
